# scripts/wled-server.service
[Unit]
Description=WLED Cube Animation Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/wled
ExecStart=/home/pi/wled/.venv/bin/uvicorn server.server:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
# Genug Zeit für mDNS-Auflösung beim Start
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
