# scripts/wled-hotkeys.service
[Unit]
Description=WLED Cube Hotkey-Daemon
After=network.target wled-server.service
Wants=wled-server.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/wled
ExecStart=/home/pi/wled/.venv/bin/python scripts/hotkeys.py --server http://localhost:8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
